# etc/beagle-server.service
# Copyright (c) 2026 Douglas P. Kingston III. MIT License - see LICENSE.
[Unit]
Description=Beagle Radio Direction Finding Server
Documentation=https://github.com/example/Beagle
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=beagle
Group=beagle

# Adjust path if installed via pip in a venv
ExecStart=/home/beagle/beagle_rdf/env/bin/beagle-server --config /etc/beagle/server.yaml

# Restart on any failure; wait 5 s between attempts.
Restart=on-failure
RestartSec=5s

# Give the process 10 s to shut down gracefully on SIGTERM
TimeoutStopSec=10s

# Keep stdout/stderr in the journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beagle-server

# Environment
Environment=PYTHONUNBUFFERED=1
# EnvironmentFile=/etc/beagle/secrets.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Server needs write access to data/ for SQLite databases
ReadWritePaths=/var/lib/beagle /tmp

[Install]
WantedBy=multi-user.target
